# etc/cron.d/desktop
03 * * * * {USER} cd /opt/hs-desktop-tool/ ; DISPLAY=:0 ~/.cabal/bin/hs-desktop-tool >> logs/refresh.log 2>&1
